# compression/systemd/headroom-proxy.service
[Unit]
Description=Headroom Compression Proxy
After=network.target

[Service]
Type=simple
ExecStart=/home/cheta/code/claude-code-proxy/.venv/bin/headroom proxy --port 8787 --mode token_headroom --openai-api-url https://openrouter.ai/api/v1 --backend openrouter --no-telemetry
Restart=on-failure
RestartSec=5
Environment="PATH=/home/cheta/code/claude-code-proxy/.venv/bin:/home/cheta/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="CUDA_VISIBLE_DEVICES=0"

[Install]
WantedBy=default.target
