# compression/systemd/surface-audit.service
[Unit]
Description=Surface system health audit

[Service]
Type=oneshot
ExecStart=%h/.local/bin/surface-audit.sh
StandardOutput=journal
StandardError=journal
